# nano/gpt_ufs_mfsroot/etc/rc.local
########################################
#
# prepare named configuration
# 
cp -fr /cfg/named/* /var/named/etc/namedb/
/etc/rc.d/named onestart
